# overlay/usb-force-host.service
[Unit]
Description=Force USB OTG host mode on rk3562
After=local-fs.target systemd-udevd.service
Wants=systemd-udevd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-force-host.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
